# etc/rc.d/handbook.in
#!/bin/sh

# PROVIDE: handbook
# REQUIRE: LOGIN NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="handbook"
rcvar="handbook_enable"

load_rc_config "$name"

: ${handbook_enable:="NO"}
: ${handbook_dir:="%%APPDIR%%"}
: ${handbook_user:="www"}
: ${handbook_bind:="127.0.0.1"}
: ${handbook_port:="9292"}
: ${handbook_env:="production"}
: ${handbook_bundle:="%%PREFIX%%/bin/bundle"}
: ${handbook_path:="%%PREFIX%%/bin:/sbin:/bin:/usr/sbin:/usr/bin"}

pidfile="/var/run/${name}.pid"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

handbook_start()
{
	if [ -f "$pidfile" ]; then
		echo "${name} is already running."
		return 1
	fi

	if [ ! -d "$handbook_dir" ]; then
		warn "${handbook_dir} does not exist"
		return 1
	fi

	echo "Starting ${name}."
	cd "$handbook_dir" || return 1
	export BUNDLE_GEMFILE="${handbook_dir}/Gemfile"
	export RACK_ENV="$handbook_env"
	export PATH="$handbook_path"
	/usr/sbin/daemon -f -p "$pidfile" -t "$name" -u "$handbook_user" \
		"$handbook_bundle" exec rackup \
		-o "$handbook_bind" \
		-p "$handbook_port" \
		config.ru
}

handbook_stop()
{
	if [ ! -f "$pidfile" ]; then
		echo "${name} is not running."
		return 1
	fi

	echo "Stopping ${name}."
	kill "$(cat "$pidfile")"
	rm -f "$pidfile"
}

handbook_status()
{
	if [ -f "$pidfile" ] && kill -0 "$(cat "$pidfile")" 2>/dev/null; then
		echo "${name} is running as pid $(cat "$pidfile")."
	else
		echo "${name} is not running."
		return 1
	fi
}

run_rc_command "$1"
